# config/files/common/etc/profile.d/run-idwtn-script.sh
#!/usr/bin/bash

# I Don't Want To Network (IDWTN)

sudo /usr/libexec/idwtn-script >/dev/null
