# config/files/common/usr/etc/profile.d/run-idwtn-disable-or-enable.sh
#!/usr/bin/bash
sudo /usr/libexec/idwtn-disable-or-enable
