# files/systemd/system/idwt-system.service
[Unit]
Description=IDWT Apply (System)

[Service]
Type=oneshot
ExecStart=/usr/bin/idwt apply-system

[Install]
WantedBy=default.target
