# files/systemd/system/idwt-apply.timer
[Unit]
Description=Scheduled IDWT Apply (every minute)

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=multi-user.target
