# config/systemd/system/idwt.service
[Unit]
Description=IDWTN service

[Service]
Type=oneshot
ExecStart=sudo /usr/libexec/idwt-script.nu

[Install]
WantedBy=default.target