# config/systemd/system/idwtn.service
[Unit]
Description=IDWTN service

[Service]
Type=oneshot
ExecStart=sudo /usr/libexec/idwtn-disable-or-enable
