# config/systemd/system/idwt.path
[Unit]
Description=Run IDWT service when config changes

[Path]
PathModified=/etc/idwt/config.json

[Install]
WantedBy=multi-user.target