# config/systemd/user/chezmoi-auto-apply.path
[Unit]
Description=Automatically apply chezmoi when directory changes

[Path]
PathModified=$HOME/.local/share/chezmoi/

[Install]
WantedBy=default.target