# config/systemd/user/idwt.service
[Unit]
Description=IDWT service

[Service]
Type=oneshot
ExecStart=/usr/libexec/idwt-apply-as-user

[Install]
WantedBy=default.target