# files/systemd/user/idwt-user.service
[Unit]
Description=IDWT Apply (User)

[Service]
Type=oneshot
ExecStart=/usr/bin/idwt apply-user

[Install]
WantedBy=default.target
