# files/systemd/user/idwt-user.path
[Unit]
Description=Run IDWT Apply (User) when config changes

[Path]
PathModified=/etc/idwt/config.yml

[Install]
WantedBy=multi-user.target
