# modules/homefiles/bluebuild-homefiles-apply.service
[Unit]
Description=Homefiles Apply

[Service]
ExecStart=/usr/bin/bluebuild-homefiles-apply.nu --mode skip
Type=oneshot

[Install]
WantedBy=default.target