# modules/homefiles/bluebuild-homefiles-apply.path
[Unit]
Description=Automatically apply homefiles when they change

[Path]
PathModified=/usr/share/bluebuild/homefiles/

[Install]
WantedBy=default.target