# config/image-pinned-etcs/common/nushell/profile.d/1-paths.nu
use std *

path add '/home/linuxbrew/.linuxbrew/bin'