# config/image-pinned-etcs/common/profile.d/alias-fastfetch.sh
#!/usr/bin/sh
alias fastfetch="/usr/bin/fastfetch -c /usr/share/ublue-os/ublue-os.jsonc --logo none"
